# infra/systemd/portfolio-green.service
[Unit]
Description=Portfolio Spring Boot (green)
After=network.target mariadb.service
Wants=network.target mariadb.service

[Service]
Type=simple
User=portfolio
Group=portfolio
WorkingDirectory=/opt/portfolio/current-green
EnvironmentFile=-/etc/portfolio/portfolio-common.env
EnvironmentFile=-/etc/portfolio/portfolio-green.env
ExecStart=/usr/bin/java -jar /opt/portfolio/current-green/portfolio.jar
Restart=always
RestartSec=5
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
